# systemd/hail-outreach.timer
[Unit]
Description=Hail Outreach Weekly Timer — fires every Monday at 12:00 local time
Requires=hail-outreach.service

[Timer]
OnCalendar=Mon *-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
